# Real_MFA/config/real_mfa_celery.service
[Unit]
Description=Real_MFA Celery Background Task Worker
After=network.target redis-server.service

[Service]
Type=simple
User=realuser
Group=www-data
WorkingDirectory=/home/realuser/real_mfa/Real_MFA
Environment="PATH=/home/realuser/real_mfa/venv/bin"
EnvironmentFile=/home/realuser/real_mfa/Real_MFA/.env

# Celery worker command
ExecStart=/home/realuser/real_mfa/venv/bin/celery -A Real_MFA worker \
    --loglevel=info \
    --logfile=/var/log/real_mfa_celery_worker.log \
    --concurrency=4

# Process management
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
